# deploy/audit-platform.service
[Unit]
Description=Audit Platform Next.js Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/audit-platform
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000
StandardOutput=append:/var/log/audit-platform/app.log
StandardError=append:/var/log/audit-platform/app-error.log

[Install]
WantedBy=multi-user.target