# deploy/audit-websocket.service
[Unit]
Description=Audit Platform WebSocket Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/audit-platform
ExecStart=/usr/bin/npx tsx ws-server.ts
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=WEBSOCKET_PORT=3001
StandardOutput=append:/var/log/audit-platform/websocket.log
StandardError=append:/var/log/audit-platform/websocket-error.log

[Install]
WantedBy=multi-user.target